# geth.service
[Unit]
Description=Go Ethereum Client
After=network.target
Wants=network.target

[Service]
User=goeth
Group=goeth
Type=simple
Restart=always
RestartSec=5
ExecStart=geth --goerli --http --datadir /var/lib/goethereum --cache 2048 --maxpeers 30 --metrics --pprof --bootnodes "enode://ef844763d501257ec2ba64a691919bd714c6c696a369fc46405b1860709b37924a32c89388f75a44e945d13ec3a53c1d06700103c654e501d0ade33c3b207018@185.181.60.210:30303,enode://46add44b9f13965f7b9875ac6b85f016f341012d84f975377573800a863526f4da19ae2c620ec73d11591fa9510e992ecc03ad0751f53cc02f7c7ed6d55c7291@94.237.54.114:30313,enode://703e48a319e32f53655a8b17d6ed1f856fb6d3f4c0892de5204912f45c80a87df1d6d3a950404137872b4b15bb56b585e5db0692a8b1ac215f7a227df117564d@185.63.173.35:30303,enode://3a7bd16841821d1951105b2fe9d6eabddc0fb77785dd60b9ef682e16917c80b808df9763def66a5b7843bcf5cdd09f17609edfeac071761a023ffb1327842137@94.237.50.17:30303"

[Install]
WantedBy=default.target